# files/etc/supervisor/conf.d/gearmand.conf
; vim:set ft=dosini:
[program:gearmand]
command = /usr/sbin/gearmand
	--log-file /var/log/libravatar/workers.log

killasgroup=true
stopasgroup=true
